# etc/profile.d/cursor.sh
export XCURSOR_THEME=Vanilla-DMZ
export XCURSOR_SIZE=48
