# etc/profile.d/unzip.sh
export UNZIP="-O CP936"
export ZIPINFO="-O CP936"
